# crontab
# run every day at 12
# need to make this every weekday at
0 16 * * 1-5 /home/mgreenway/bullet/get_bullets.py -r &> /dev/null
